# root/etc/crontabs/abc
# min   hour    day     month   weekday command
*       *       *       *       *       cd /app/www; /usr/bin/php artisan schedule:run 2>&1
